# icinga2.ini
[program:icinga2]
command=/usr/sbin/icinga2 daemon -e /var/log/icinga2/error.log
